# packaging/owllock.service
[Unit]
Description=OwlLock background service
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/home/snowowl/Projects/faceLock/scripts/owllock-launcher.sh --background
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
