# packaging/owllock.service
[Unit]
Description=OwlLock background service
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart={launcher} --background
WorkingDirectory={workdir}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
